# USB_mount/rc.local
#!/bin/sh
/home/pi/kumikomi/USB_mount/usb_monitor.py